# launchd/com.bazaarskills.chrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bazaarskills.chrome</string>

  <!-- TEMPLATE: launchd/install_daemon.sh substitutes __RUNTIME__ (this checkout) at install. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__RUNTIME__/bin/chrome_debug.sh</string>
  </array>

  <!-- start at login, restart if Chrome is closed/crashes (always-on, dedicated profile) -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <!-- give it a moment between restarts so it doesn't hot-loop on a bad launch -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__RUNTIME__/logs/chrome.out.log</string>
  <key>StandardErrorPath</key>
  <string>__RUNTIME__/logs/chrome.err.log</string>
</dict>
</plist>
