# launchd/com.bazaarskills.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bazaarskills.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__RUNTIME__/bin/agent_daemon.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__RUNTIME__</string>

  <!-- This is a TEMPLATE. launchd/install_daemon.sh substitutes __RUNTIME__ (this checkout) and
       __PATH__ (the dirs where your claude / npx / node / curl live, detected at install) before
       copying the plist into ~/Library/LaunchAgents. launchd jobs otherwise get a minimal PATH. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <!-- Pin the agent runtime so the daemon never autodetects `codex` (which the runner refuses
         → rc=3, a whole pass lost). Only claude-code is wired today; see harness_run.py. -->
    <key>BAZAAR_HARNESS</key>
    <string>claude-code</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__RUNTIME__/logs/agent.out.log</string>
  <key>StandardErrorPath</key>
  <string>__RUNTIME__/logs/agent.err.log</string>
</dict>
</plist>
